# deploy/systemd/system/kitt-jmri.service

[Unit]
Description=KITT JMRI (Headless PanelPro)
After=network-online.target mosquitto.service
Wants=network-online.target
PartOf=kitt.target

# Optional: if your EX-CommandStation is reachable only after it boots on LAN,
# keep network-online to ensure DHCP/links are ready.

[Service]
Type=simple
User=pi
Group=pi
# Set a stable working dir where JMRI is installed
WorkingDirectory=/opt/jmri

# Headless rendering flag for Java-based apps
Environment=JAVA_TOOL_OPTIONS=-Djava.awt.headless=true

# Optional: KITT env (e.g., broker host if not localhost)
EnvironmentFile=/etc/kitt/kitt.env

# ---- ExecStart options ----
# Choose ONE of the following ExecStart forms:

# (A) If you have the PanelPro launcher script in /opt/jmri:
# ExecStart=/opt/jmri/PanelPro --no-splash --profile KITT

# (B) Generic Java invocation (if you know the PanelPro jar path):
# ExecStart=/usr/bin/java -Djava.awt.headless=true -jar /opt/jmri/lib/PanelPro.jar --profile KITT

# (C) Wrapper script (recommended if your paths/settings vary):
# ExecStart=/opt/kitt/src/jmri/run_headless.sh

# --- pick A, B, or C, then remove the others. ---

Restart=on-failure
RestartSec=3
# Hardening (tune as needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=kitt.target
